# k8s/cluster/k8s/roles/zcm-config/files/zcm-config.ini
[program:zcm-config]
command=/bin/bash -c '. /zcm/supervisor/zcm-env.profile; \
    [ `docker ps -a -f "name=zcm-config"| wc -l` -gt 1 ] && \
    exec docker start -a zcm-config || \
        docker run --name zcm-config \
        -e SERVER_IP=${ZCM_SERVER_LIST} \
        --net=host \
        -v /zcm/config:/zcm/config \
        -e CLOUD_APP_NAME=paas_zcm-config \
        -e ZCM_IMAGE_PORT=8280 \
        -v /etc/localtime:/etc/localtime:ro \
        -v /etc/timezone:/etc/timezone:ro \
        -v /tmp/zlogs:/tmp/zlogs \
        -v /tmp/zcore:/tmp/zcore \
        -d \
        ${ZCM_HARBOR}/zcm9/zcm-git:${ZCM_CONFIG_TAG}; \
    exec docker start -a zcm-config'
process_name= zcm-config
numprocs=1
autostart=true
autorestart=true
startretries = 6
user=root
stdout_logfile=/var/log/zcm/zcm-config.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/zcm/zcm-config.log
stderr_logfile_maxbytes=10MB
stopwaitsecs=20
