# k8s/cluster/k8s/roles/install_proxy/files/zcm-node-proxy.ini
[program:zcm-node-proxy]
command=/bin/bash -c \
    '. /zcm/supervisor/zcm-env.profile; \
    [ `docker ps -a -f "name=zcm-node-proxy"| wc -l` -gt 1 ] && \
    exec docker start -a zcm-node-proxy || \
    docker run --name zcm-node-proxy \
        -v /zcm/node-proxy/nginx.conf:/usr/local/openresty/nginx/conf/nginx.conf \
        --restart=on-failure:5 \
        --memory=512M \
        --net=host \
        -e CLOUD_APP_NAME=paas_zcm-node-proxy \
        -v /etc/localtime:/etc/localtime:ro \
        -v /etc/timezone:/etc/timezone:ro \
        -v /tmp/zlogs:/tmp/zlogs \
        -v /tmp/zcore:/tmp/zcore \
        -d \
        ${ZCM_HARBOR}/${ZCM_HARBOR_PROJECT}/zcm-openresty:${ZCM_NODE_PROXY_TAG}; \
    exec docker start -a zcm-node-proxy'
process_name= zcm-node-proxy
numprocs=1
autostart=true
autorestart=true
startretries = 6
user=root
stdout_logfile=/var/log/zcm/zcm-node-proxy.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/zcm/zcm-node-proxy.log
stderr_logfile_maxbytes=10MB
stopwaitsecs=20

